# deploy/systemd/senior-pomidor-backup@.service
[Unit]
Description=Senior Pomidor %i backup
Requires=docker.service
After=docker.service senior-pomidor.service

[Service]
Type=oneshot
User=root
EnvironmentFile=/srv/secret/senior-pomidor.env
ExecStart=/srv/automation/backup.sh %i
